# clients/crontab
0	*	*	*	*	$HOME/hephaestus/trunk/clients/hephaestus > /dev/null
